# enhanceposition_tata/enhanceposition/adrd/adrd.service
[Unit]
Description=ADR Daemon

[Service]
ExecStart=/usr/bin/adrd
Restart=always
Requires=mnld

[Install]
WantedBy=basic.target
