# gps_tata/gps/mtk_mnld/mnld.service
[Unit]
Description=MNL Daemon

[Service]
ExecStart=/usr/bin/mnld0
Restart=always
#StandardOutput=kmsg+console

[Install]
WantedBy=basic.target
